# roles/app_redis/templates/redis.service.j2
[Unit]
Description=redis service {{ username }}

[Service]
Restart=always
ExecStart={{ redis_server_binary }} /home/{{ username }}/{{ config_filename }}

[Install]
WantedBy=default.target
