# roles/server_borg/templates/remote_storage.automount.j2
[Unit]
Description=BackupSpace

[Automount]
Where={{ borg_remote_storage_mount_path }}
DirectoryMode=0700
TimeoutIdleSec=30

[Install]
WantedBy=multi-user.target
